# infra/bootstrap/units/pwfg-proxy.service
[Unit]
Description=pwfg LLM brokering proxy (hides the Anthropic key, caps cost, audits)
After=network-online.target
Wants=network-online.target
# The loop must not start until the proxy the agent talks to is up.
Before=pwfg-loop.service

[Service]
User=proxy
Group=pwfgkey
WorkingDirectory=/srv/pwfg
# The real key arrives via the kernel keyring, NOT the environment or /proc.
# Provision the SecureString to this path out of band (see P1-provisioning.md).
LoadCredential=anthropic_key:/run/pwfg/anthropic_key
Environment=PWFG_PROXY_HOST=127.0.0.1
Environment=PWFG_PROXY_PORT=8787
Environment=PWFG_PROXY_STATE=/srv/pwfg/proxy
Environment=PWFG_CONTROL_DIR=/srv/pwfg/control
Environment=PWFG_PROXY_UPSTREAM=https://api.anthropic.com
# Reach the upstream THROUGH the in-VPC Squid forward proxy (domain allowlist). The
# agent-host subnet is isolated; Squid is the only internet path. The real Squid
# private IP is the PwfgEgress `SquidPrivateIp` CfnOutput — fill it OUT OF BAND (a
# systemd drop-in, like the key), since cloud-init is static text and a CDK token
# cannot be interpolated post-read. The placeholder fails closed (an unresolved host
# stalls the loop rather than leaking); see docs/P1-provisioning.md.
Environment=PWFG_PROXY_FORWARD=http://SQUID_PRIVATE_IP:3128
# Hard ceiling — the authoritative cost cap. Tune per run.
Environment=PWFG_PROXY_MAX_COST_USD=25.00
# Pin httpx <1 so the AsyncClient(proxy=...) chaining stays API-stable at runtime.
ExecStart=/usr/bin/uv run --python 3.13 --with starlette --with 'httpx>=0.28,<1' --with uvicorn python -m proxy.app
WorkingDirectory=/srv/pwfg/proxy-src
Restart=on-failure
RestartSec=2

# Hardening: no core dumps (key in memory), no new privileges, private tmp,
# read-only system, write only the proxy state + control dirs.
LimitCORE=0
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/srv/pwfg/proxy /srv/pwfg/control
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
